# etc/zephyr-shell.service
[Unit]
Description=Zephyr shell UART over TCP (port 4444)
After=network.target

[Service]
ExecStart=/usr/bin/socat TCP-LISTEN:4444,reuseaddr,fork /dev/ttyYahboom1,b921600,raw,echo=0
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
